# docker_files/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=Xvfb :99 -screen 0 1024x768x16

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"